# deploy/vigil-collector.service
[Unit]
Description=Vigil Collector - Server Monitor (v2)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /root/server-monitor/vigil/main.py
WorkingDirectory=/root/server-monitor/vigil
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
